# build/queue/yii-queue.ini
[program:yii-queue]
command=/data/www/oms/console/yii queue/listen
stopsignal=QUIT
autorestart=true
autostart=true
