# build/nginx/nginx.ini
[program:nginx]
command=/usr/local/nginx/sbin/nginx -g "daemon off;" -c /usr/local/nginx/conf/nginx.conf
stopsignal=QUIT
autorestart=true
autostart=true